# bin/udev/95-hdmi-plug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/mathias/bin/udev/hdmi-plug.sh"
